# pi_home/crontab
0 0 * * * find /home/pi/Monitor -mtime 2 -delete
